# lab_pi_session_poller.service
[Unit]
Description=Lab Pi Session Poller
After=network.target

[Service]
Type=simple
User=abhi
WorkingDirectory=/home/abhi/lab-pi
EnvironmentFile=/home/abhi/lab-pi/.env
ExecStart=/usr/bin/python3 /home/abhi/lab-pi/lab_pi_session_poller.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
