# systemd/audio_stream.service
[Unit]
Description=Virtual Lab - Audio Streaming Service (Lab-Pi)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=dese
WorkingDirectory=/home/dese/lab-pi/lab-pi/Audio
ExecStart=/usr/bin/python3 /home/dese/lab-pi/lab-pi/Audio/server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1

# Hardening
ProtectSystem=full
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
