# Audio/services/audio_stream.service
[Unit]
Description=Virtual Lab - Audio Capture Service (Lab-Pi)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/lab-pi
ExecStart=/usr/bin/python3 %h/lab-pi/audio_capture.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=VLAB_PI_ID=LabPi1
Environment=MASTER_URL=http://192.168.1.5:5000

# Hardening
ProtectSystem=full
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
