# Audio/services/audio_stream.service
[Unit]
Description=Virtual Lab - Audio Streaming Service (Lab-Pi)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/lab-pi/sound
ExecStart=/usr/bin/python3 %h/lab-pi/sound/server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1

# Hardening
ProtectSystem=full
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
