# Audio/services/mjpg-streamer.service
[Unit]
Description=Virtual Lab - MJPG/USTREAMER Video Stream Service
After=network.target

[Service]
Type=simple
User=%i
ExecStart=/usr/local/bin/ustreamer --host 0.0.0.0 --port 8080 --device /dev/video0 --format jpeg --resolution 640x480 --fps 15
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
